# devices/xiaomi_ax3600/files.edma/etc/uci-defaults/998-net-offload
uci -q batch << EOI
# EDMA (CPU-bound) data path: enable RPS packet steering and software flow offload to
# claw back NAT throughput. IPQ807x EDMA has no hardware offload, so leave _hw off.
set network.globals.packet_steering='1'
set network.@device[0].flow_offloading='1'
set network.@device[0].flow_offloading_hw='0'
EOI
