# devices/xiaomi_ax3600/files.edma-nss/etc/uci-defaults/998-net-offload
uci -q batch << EOI
# NSS data path: kernel packet steering and flow offload stay off — the NSS cores and
# ECM handle acceleration. (Rest of the network defaults live in 999-QOL_config.)
set network.globals.packet_steering='0'
set network.@device[0].flow_offloading='0'
set network.@device[0].flow_offloading_hw='0'
EOI

# Enable the NSS bring-up service. Runs on every first boot after a flash, so
# the rc.d link always exists; the durable opt-out is the uci flag the service
# checks (nss.general.enabled='0'), which survives sysupgrade with the config.
/etc/init.d/nss enable
