# devices/xiaomi_ax3600/files.edma-nss/etc/uci-defaults/999-sqm-settings
#!/bin/sh
# Provision the NSS SQM (nss-edma.qos) shaper. Create a proper sqm.wan section
# only if absent (so user-tuned rates survive sysupgrade) and remove the legacy
# bogus sqm.eth1 section (it was named 'eth1' but pointed at interface 'wan' and
# carried no download/upload rates, so it shaped nothing and produced a spurious
# second SQM start at boot). Rates below are for a ~300/150 Mbit line; adjust to
# your link.
if [ -z "$(uci -q get sqm.wan)" ]; then
	uci -q batch <<-EOI
		set sqm.wan=queue
		set sqm.wan.enabled='1'
		set sqm.wan.interface='wan'
		set sqm.wan.download='285000'
		set sqm.wan.upload='142500'
		set sqm.wan.qdisc='fq_codel'
		set sqm.wan.script='nss-edma.qos'
		set sqm.wan.qdisc_advanced='1'
		set sqm.wan.qdisc_really_really_advanced='1'
		set sqm.wan.itarget='4ms'
		set sqm.wan.etarget='4ms'
		set sqm.wan.iqdisc_opts='interval 50ms quantum 304'
		set sqm.wan.eqdisc_opts='interval 50ms quantum 304'
	EOI
fi
uci -q delete sqm.eth1
uci -q commit sqm
exit 0
