# files/etc/uci-defaults/999-QOL_config
uci -q batch << EOI
uci set pbuf.opt.memory_profile='auto'
uci set wireless.default_radio0.dtim_period='1'
uci set wireless.default_radio0.max_inactivity='86400'
uci set wireless.default_radio0.disassoc_low_ack='0'
uci set wireless.default_radio0.wnm_sleep_mode='1'
uci set wireless.radio1.country='US'
uci set wireless.default_radio1.dtim_period='1'
uci set wireless.default_radio1.max_inactivity='86400'
uci set wireless.default_radio1.disassoc_low_ack='0'
uci set wireless.radio1.htmode='HE160'
uci set wireless.radio1.channel='44'
uci set wireless.radio0.country='US'
uci set wireless.default_radio1.wnm_sleep_mode='1'
uci set wireless.default_radio2.dtim_period='1'
uci set wireless.default_radio2.max_inactivity='86400'
uci set wireless.default_radio2.disassoc_low_ack='0'
uci set wireless.default_radio2.wnm_sleep_mode='1'
uci set wireless.radio2.country='US'
uci set firewall.@defaults[0].flow_offloading=0
uci set network.globals.packet_steering=0
uci set ecm.@general[0].enable_bridge_filtering=0
EOI
