# files/etc/uci-defaults/999-QOL_config
uci -q batch << EOI
set pbuf.opt.memory_profile='auto'
set wireless.default_radio0.dtim_period='1'
set wireless.default_radio0.max_inactivity='86400'
set wireless.default_radio0.disassoc_low_ack='0'
set wireless.default_radio0.wnm_sleep_mode='1'
set wireless.radio0.htmode='HE160'
set wireless.radio0.channel='44'
set wireless.radio0.country='US'
set wireless.default_radio1.dtim_period='1'
set wireless.default_radio1.max_inactivity='86400'
set wireless.default_radio1.disassoc_low_ack='0'
set wireless.default_radio1.wnm_sleep_mode='1'
set wireless.radio1.country='US'
set firewall.@defaults[0].flow_offloading='0'
set network.globals.packet_steering='0'
set ecm.@general[0].enable_bridge_filtering='0'
set network.@device[0].igmp_snooping='1'
set network.@device[0].multicast_querier='1'
set firewall.@zone[1].input='DROP'
set firewall.@zone[1].forward='DROP'
set uhttpd.main.redirect_https='1'
EOI
