# files/etc/uci-defaults/999-QOL_config
uci -q batch << EOI
# Set pbuf options
set pbuf.opt.auto_scale='on'
set pbuf.opt.memory_profile='1gb'
set pbuf.opt.scaling_governor='performance'

# Wireless default settings for radio0
set wireless.default_radio0.dtim_period='1'
set wireless.default_radio0.max_inactivity='86400'
set wireless.default_radio0.disassoc_low_ack='0'
set wireless.default_radio0.wnm_sleep_mode='1'
set wireless.default_radio0.ieee80211k='1'
set wireless.default_radio0.time_advertisement='2'
set wireless.default_radio0.bss_transition='1'
set wireless.default_radio0.proxy_arp='1'

# Wireless radio0 specific settings
set wireless.radio0.htmode='HE160'
set wireless.radio0.channel='auto'
set wireless.radio0.country='DE'

# Wireless default settings for radio1
set wireless.default_radio1.dtim_period='1'
set wireless.default_radio1.max_inactivity='86400'
set wireless.default_radio1.disassoc_low_ack='0'
set wireless.default_radio1.wnm_sleep_mode='1'
set wireless.default_radio1.ieee80211k='1'
set wireless.default_radio1.time_advertisement='2'
set wireless.default_radio1.bss_transition='1'
set wireless.default_radio1.proxy_arp='1'

# Wireless radio1 specific settings
set wireless.radio1.htmode='HE40'
set wireless.radio1.channel='auto'
set wireless.radio1.country='DE'

# Firewall settings
set firewall.@defaults[0].flow_offloading='0'
set firewall.@zone[1].input='DROP'
set firewall.@zone[1].forward='DROP'

# Network settings
set network.globals.packet_steering='0'
set network.@device[0].igmp_snooping='1'
set network.@device[0].multicast_querier='1'

# ECM settings
set ecm.@general[0].enable_bridge_filtering='0'

# uHTTPd settings
set uhttpd.main.redirect_https='1'

# BanIP settings
set banip.global.ban_cores='1'
set banip.global.ban_reportelements='0'
set banip.global.ban_splitsize='1024'
set banip.global.ban_blockinput='cinsscore' 'debl' 'turris'
set banip.global.ban_blockforwardwan='cinsscore' 'debl' 'turris'
set banip.global.ban_blockforwardlan='doh'
set banip.global.ban_countries='BY','CD','CF','CU','ET','HK','IQ','IR','KP','LB','LY','MM','RU','SD','SY','VE','YE','ZW'

# Disable IPv6
set network.lan.ipv6='0'
set network.wan.ipv6='0'
network.@device[0].ipv6='0'
set dhcp.lan.dhcpv6=disabled
-q delete dhcp.lan.dhcpv6
-q delete dhcp.lan.ra
set network.lan.delegate='0'
EOI
