# devices/xiaomi_ax3600/files/etc/uci-defaults/999-QOL_config
uci -q batch << EOI

# Wireless defaults live in the static /etc/config/wireless (radios enabled
# and configured from the first boot - no detection race).

# Firewall settings
set firewall.@zone[1].input='DROP'
set firewall.@zone[1].forward='DROP'

# Network settings (BCP38 anti-spoofing on WAN).
# Packet steering / flow offload are variant-specific — see files.<variant>/etc/uci-defaults/998-net-offload.
set bcp38.cfg018e45.enabled='1'
set bcp38.cfg018e45.interface='wan'

# NTP: serve time locally (sysntpd) and advertise this router to LAN clients
# via DHCP option 42 (which carries IP addresses, not hostnames). The address
# is derived from the live LAN config, not hardcoded.
set system.ntp.enable_server='1'

# uHTTPd settings
set uhttpd.main.redirect_https='1'

# OpenSSL OQS provider
set openssl.oqsprovider=provider
set openssl.oqsprovider.enabled='1'

EOI

# DHCP option 42 needs the router's actual LAN address - derive it from the
# live config (uci-defaults run after the board defaults exist).
lan_ip=$(uci -q get network.lan.ipaddr)
[ -n "$lan_ip" ] && uci -q add_list dhcp.lan.dhcp_option="42,${lan_ip%%/*}"
exit 0
