# devices/xiaomi_ax3600/files.edma-nss/etc/init.d/nss
#!/bin/sh /etc/rc.common
# Bring up the NSS offload data plane after boot by running /usr/sbin/nss-up
# (firmware arm, bridge/multicast/ECM/PPPoE managers, Wi-Fi offload flip, SQM).
# The plane cannot be detached at runtime; a reboot always returns the device
# to the plain host stack. To stay on the host stack permanently:
#   uci set nss.general.enabled='0'; uci commit nss
# (the uci flag survives sysupgrade, unlike an rc.d disable).

START=99
USE_PROCD=1

start_service() {
	[ "$(uci -q get nss.general.enabled)" = "0" ] && {
		logger -t nss "NSS bring-up disabled via uci (nss.general.enabled=0), staying on the host stack"
		return 0
	}
	procd_open_instance nss-up
	procd_set_param command /usr/sbin/nss-up
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_close_instance
}

stop_service() {
	echo "The NSS plane cannot be detached at runtime; reboot to return to the host stack."
	echo "To keep it off across reboots: uci set nss.general.enabled='0'; uci commit nss"
}
